# com.chadhomes.dashboard-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job for the dashboard watcher (post-restructure).

  Install:
    cp com.chadhomes.dashboard-watcher.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.chadhomes.dashboard-watcher.plist

  Stop:
    launchctl unload ~/Library/LaunchAgents/com.chadhomes.dashboard-watcher.plist

  Reload after edits:
    launchctl unload ~/Library/LaunchAgents/com.chadhomes.dashboard-watcher.plist
    launchctl load   ~/Library/LaunchAgents/com.chadhomes.dashboard-watcher.plist

  Inspect:
    launchctl list | grep chadhomes
    tail -f ~/Projects/home-builder-agent/watcher.log
    tail -f /tmp/dashboard-watcher.stderr.log

  Post-restructure changes:
    - ProgramArguments now invokes the package as a module (`python3 -m
      home_builder_agent.watchers.dashboard`) rather than a flat-file path.
    - WorkingDirectory still ~/Projects/home-builder-agent/ — that's where
      .env, credentials.json, token.json, and .watcher_state.json live.
    - PYTHONPATH set to project root so the package import resolves without
      requiring `pip install -e .` first (belt-and-suspenders; the editable
      install also makes this work).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chadhomes.dashboard-watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
    <string>-m</string>
    <string>home_builder_agent.watchers.dashboard</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/connorpatton/Projects/home-builder-agent</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/Users/connorpatton/Projects/home-builder-agent</string>
  </dict>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/dashboard-watcher.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/dashboard-watcher.stderr.log</string>
</dict>
</plist>
