# com.chadhomes.tracker-mirror.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job for the Postgres → Google Sheets Tracker mirror worker.

  Per ADR 2026-05-11 ("Postgres as Canonical State; Google Sheets as
  Read-Only Mirror"), this worker overwrites the Master Schedule and
  Project Info tabs of every active project's Tracker from Postgres
  every 5 minutes. The Tracker is a *display* layer; Postgres is
  canonical.

  DO NOT INSTALL until CTO has reviewed the mirror_worker code path.
  This file is intentionally checked in but unloaded; Connor loads it
  manually once the worker passes a dry-run smoke test.

  Install (after CTO sign-off):
    cp com.chadhomes.tracker-mirror.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.chadhomes.tracker-mirror.plist

  Stop:
    launchctl unload ~/Library/LaunchAgents/com.chadhomes.tracker-mirror.plist

  Reload after edits:
    launchctl unload ~/Library/LaunchAgents/com.chadhomes.tracker-mirror.plist
    launchctl load   ~/Library/LaunchAgents/com.chadhomes.tracker-mirror.plist

  Inspect:
    launchctl list | grep chadhomes
    tail -f /tmp/tracker-mirror.stderr.log
    tail -f /tmp/tracker-mirror.stdout.log

  Pattern source: com.chadhomes.dashboard-watcher.plist (same env-bootstrap,
  same log conventions, longer StartInterval).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chadhomes.tracker-mirror</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
    <string>-m</string>
    <string>home_builder_agent.scheduling.mirror_worker</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/connorpatton/Projects/home-builder-agent</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/Users/connorpatton/Projects/home-builder-agent</string>
  </dict>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/tracker-mirror.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/tracker-mirror.stderr.log</string>
</dict>
</plist>
</content>
</invoke>